# lib/dpkg/info/fio.list
/.
/etc
/etc/init.d
/etc/init.d/fio
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/fio.service
/usr
/usr/bin
/usr/bin/fio
/usr/bin/fio-btrace2fio
/usr/bin/fio-dedupe
/usr/bin/fio-genzipf
/usr/bin/fio-verify-state
/usr/bin/fio2gnuplot
/usr/bin/fio_generate_plots
/usr/bin/fio_jsonplus_clat2csv
/usr/bin/genfio
/usr/share
/usr/share/doc
/usr/share/doc/fio
/usr/share/doc/fio/HOWTO.rst.gz
/usr/share/doc/fio/MORAL-LICENSE
/usr/share/doc/fio/NEWS.Debian.gz
/usr/share/doc/fio/README.Debian
/usr/share/doc/fio/README.rst.gz
/usr/share/doc/fio/REPORTING-BUGS
/usr/share/doc/fio/changelog.Debian.gz
/usr/share/doc/fio/copyright
/usr/share/doc/fio/fio-histo-log-pctiles.py.gz
/usr/share/doc/fio/fiograph.conf
/usr/share/doc/fio/fiograph.py.gz
/usr/share/doc/fio/fiologparser.py.gz
/usr/share/doc/fio/fiologparser_hist.py.1.gz
/usr/share/doc/fio/fiologparser_hist.py.gz
/usr/share/doc/fio/half-bins.py
/usr/share/doc-base
/usr/share/doc-base/fio.fio
/usr/share/fio
/usr/share/fio/graph2D.gpm
/usr/share/fio/graph3D.gpm
/usr/share/fio/math.gpm
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/fio.1.gz
/usr/share/man/man1/fio2gnuplot.1.gz
/usr/share/man/man1/fio_generate_plots.1.gz
