# lib/dpkg/info/raspi-config.conffiles
remove-on-upgrade /etc/init.d/raspi-config
remove-on-upgrade /etc/default/cpu_governor
/etc/sudoers.d/010_proxy
