# lib/dpkg/info/rpi-eeprom.list
/.
/etc
/etc/default
/etc/default/rpi-eeprom-update
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rpi-eeprom-update.service
/usr
/usr/bin
/usr/bin/rpi-bootloader-key-convert
/usr/bin/rpi-eeprom-config
/usr/bin/rpi-eeprom-digest
/usr/bin/rpi-eeprom-update
/usr/bin/rpi-otp-private-key
/usr/bin/rpi-sign-bootcode
/usr/lib
/usr/lib/firmware
/usr/lib/firmware/raspberrypi
/usr/lib/firmware/raspberrypi/bootloader-2711
/usr/lib/firmware/raspberrypi/bootloader-2711/default
/usr/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2023-01-11.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2024-04-15.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2025-02-11.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/default/recovery.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/default/vl805-000138c0.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2023-01-11.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2023-05-11.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2024-01-18.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2024-01-22.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2024-04-15.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2024-04-17.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2024-05-17.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2024-07-30.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2024-09-05.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2024-10-10.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2024-10-21.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2024-12-07.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2025-02-11.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2025-02-24.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2025-05-08.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/recovery.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/latest/vl805-000138c0.bin
/usr/lib/firmware/raspberrypi/bootloader-2711/release-notes.md
/usr/lib/firmware/raspberrypi/bootloader-2712
/usr/lib/firmware/raspberrypi/bootloader-2712/default
/usr/lib/firmware/raspberrypi/bootloader-2712/default/pieeprom-2024-11-12.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/default/pieeprom-2025-01-22.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/default/pieeprom-2025-02-12.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/default/pieeprom-2025-03-10.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/default/recovery.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/latest
/usr/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2024-11-12.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2025-01-22.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2025-02-12.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2025-03-10.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2025-03-19.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2025-03-27.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2025-04-07.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2025-05-08.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/latest/recovery.bin
/usr/lib/firmware/raspberrypi/bootloader-2712/release-notes.md
/usr/share
/usr/share/doc
/usr/share/doc/rpi-eeprom
/usr/share/doc/rpi-eeprom/LICENCE.bootloader
/usr/share/doc/rpi-eeprom/changelog.Debian.gz
/usr/share/doc/rpi-eeprom/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/rpi-eeprom-config.1.gz
/usr/share/man/man1/rpi-eeprom-update.1.gz
/var
/var/lib
/var/lib/raspberrypi
/var/lib/raspberrypi/bootloader
/var/lib/raspberrypi/bootloader/backup
/usr/lib/firmware/raspberrypi/bootloader-2711/beta
/usr/lib/firmware/raspberrypi/bootloader-2711/critical
/usr/lib/firmware/raspberrypi/bootloader-2711/stable
/usr/lib/firmware/raspberrypi/bootloader-2712/beta
/usr/lib/firmware/raspberrypi/bootloader-2712/critical
/usr/lib/firmware/raspberrypi/bootloader-2712/stable
