# lib/dpkg/info/wayvnc.list
/.
/etc
/etc/wayvnc
/etc/wayvnc/config
/usr
/usr/bin
/usr/bin/wayvnc
/usr/bin/wayvncctl
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/wayvnc-control.service
/usr/lib/systemd/system/wayvnc-generate-keys.service
/usr/lib/systemd/system/wayvnc.service
/usr/lib/sysusers.d
/usr/lib/sysusers.d/vnc.conf
/usr/sbin
/usr/sbin/wayvnc-control.py
/usr/sbin/wayvnc-generate-keys.sh
/usr/sbin/wayvnc-run.sh
/usr/share
/usr/share/doc
/usr/share/doc/wayvnc
/usr/share/doc/wayvnc/changelog.Debian.gz
/usr/share/doc/wayvnc/copyright
/usr/share/doc/wayvnc/examples
/usr/share/doc/wayvnc/examples/README.md
/usr/share/doc/wayvnc/examples/auto-attach.py
/usr/share/doc/wayvnc/examples/event-watcher
/usr/share/doc/wayvnc/examples/single-output-sway
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/wayvnc.1.gz
/usr/share/man/man1/wayvncctl.1.gz
